# frp/client/frpc.service
[Unit]
Description=frpc
After=network.target
 
[Service]
ExecStart=/usr/local/frpc/start.sh
 
[Install]
WantedBy=multi-user.target

